# cron.d
@reboot (sleep 300 && updatedb) >> /var/log/updatedb.log 2>&1 &
